# supervisord.conf

[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)
childlogdir=/var/log/supervisor ; ('AUTO' child log dir, default $TEMP)

[program:avahi]
command=avahi-daemon --no-drop-root

[program:snapserver]
command=snapserver -s "spotify:///usr/local/bin/librespot?name=Spotify&devicename=%(ENV_SPOTIFY_NAME)s&bitrate=%(ENV_SPOTIFY_BITRATE)s" -s "airplay:///usr/bin/shairport-sync?name=Airplay&devicename=%(ENV_SPOTIFY_NAME)s"
;command=snapserver -s "airplay:///usr/bin/shairport-sync?name=Airplay&devicename=%(ENV_SPOTIFY_NAME)s"
;stdout_logfile=/dev/null
;stderr_logfile=/dev/null

[program:snapcastr]
command=snapcastrd --port=80 --sc_host=%(ENV_SNAPCAST_HOST)s
;stdout_logfile=/dev/null
;stderr_logfile=/dev/null
